# com.dansiwiec.synthia.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dansiwiec.synthia</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/dansiwiec/.local/bin/uv</string>
        <string>run</string>
        <string>uvicorn</string>
        <string>synthia.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8003</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/dansiwiec/workspace/synthia/logs/error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/dansiwiec/workspace/synthia/logs/stdout.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/dansiwiec/workspace/synthia</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/dansiwiec/.nvm/versions/node/v22.20.0/bin:/Users/dansiwiec/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
